# root/etc/services.d/rsyslog/run
#!/usr/bin/with-contenv bash

exec rsyslogd -n -f /etc/rsyslog.conf
